# deploy/hustle-agent.service
[Unit]
Description=Hustle Agent — Autonomous Money-Making AI
After=network.target

[Service]
Type=simple
User=hustle
WorkingDirectory=/opt/hustle-agent
ExecStart=/usr/bin/python3 /opt/hustle-agent/agent/engine.py loop --interval 300
Restart=on-failure
RestartSec=30
Environment=ANTHROPIC_API_KEY=your_key_here

# Logging handled by the agent itself (logs/ directory)
StandardOutput=append:/opt/hustle-agent/logs/daemon.log
StandardError=append:/opt/hustle-agent/logs/daemon.log

[Install]
WantedBy=multi-user.target
